# roles/docker-common/files/Clean-up-docker.service
[Unit]
Description=Clean up docker

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleanup_docker.sh
WorkingDirectory=/usr/local/bin/
  